# deploy/lumastir.service
[Unit]
Description=Lumastir API Server
After=network.target

[Service]
User=pi
WorkingDirectory=~/Projects/lumastir

# Point to the config specific to this physical device

Environment="LUMASTIR_CONFIG=~/Projects/lumastir/configs/3led_3motor.yaml"

# 'uv run' handles the venv creation and dependency syncing automatically
# Ensure uv is installed: curl -LsSf https://astral.sh/uv/install.sh | sh

ExecStart=~/.local/bin/uv run lumastir-server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

