# deploy/lumastir.service
[Unit]
Description=Lumastir API Server
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/lumastir
# Point to the config specific to this physical device
Environment="LUMASTIR_CONFIG=/home/pi/lumastir/configs/3led_6fan.yaml"
# 'uv run' handles the venv creation and dependency syncing automatically
# Ensure uv is installed: curl -LsSf https://astral.sh/uv/install.sh | sh
ExecStart=/home/pi/.local/bin/uv run lumastir-server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

